# services/teslausb-logrotate.timer
[Unit]
Description=TeslaUSB Daily Log Rotation Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
