# services/teslausb-bgpreview.service
[Unit]
Description=TeslaUSB CL Background Preview Generator
After=network-online.target teslausb-web.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/radxa_data/teslausb
ExecStart=/usr/bin/python3 -u bg_preview_generator.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/teslausb-bgpreview.log
StandardError=append:/var/log/teslausb-bgpreview.log

[Install]
WantedBy=multi-user.target
