# services/wifi-quick-check.service
[Unit]
Description=WiFi Quick Check Service (Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/radxa_data/teslausb/wifi_service.py --quick
StandardOutput=journal
StandardError=journal
