# RD-T/teslausb-boot-notify.service
[Unit]
Description=TeslaUSB Boot Notification (WeChat)
Documentation=https://github.com/mphacker/TeslaUSB
After=network-online.target teslausb-web.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/radxa_data/teslausb
ExecStart=/usr/bin/python3 /opt/radxa_data/teslausb/boot_notify.py
StandardOutput=append:/var/log/teslausb-boot-notify.log
StandardError=append:/var/log/teslausb-boot-notify.log
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
